# contrib/deb/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable dcvix-director.service 