# contrib/deb/prerm
#!/bin/sh
set -e

systemctl stop dcvix-director.service || true
systemctl disable dcvix-director.service || true 